# onifast-ftp.service
[Unit]
Description=Onifast FTP Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/root/onifast
ExecStart=/usr/local/bin/onifast-ftp
Restart=always
RestartSec=5
StandardOutput=append:/home/root/onifast/logs/ftp.log
StandardError=append:/home/root/onifast/logs/ftp.log

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
