# onifast-proxy.service
[Unit]
Description=Onifast IP Proxy Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/root/onifast
ExecStart=/usr/local/bin/onifast-proxy
Restart=always
RestartSec=5
StandardOutput=append:/var/log/onifast/proxy.log
StandardError=append:/var/log/onifast/proxy.log

[Install]
WantedBy=multi-user.target
